# my-cp-badges-server-py.service
[Unit]
Description=CP-Badges.

[Service]
Type=simple
ExecStart=/home/CP-Badges/prod-run.sh

[Install]
WantedBy=graphical.target
